# packaging/postinst
#! /bin/bash
systemctl daemon-reload
systemctl enable sbsolved
systemctl start sbsolved