# packaging/sbsolved.service
[Unit]
Description = Daemon for SpellingBeeSolver

[Service]
Type = simple
ExecStart = /usr/bin/sbsolved

[Install]
WantedBy = multi-user.target